# client/osquery/honeypot.osquery.conf
{
	"options": {
		"config_plugin": "filesystem",
		"logger_plugin": "filesystem",
		"utc": "true"
	},
	"schedule": {
		"network_processes": {
			"query": "SELECT processes.pid, processes.name as process_name, etc_protocols.name as protocol_name, listening_ports.port, listening_ports.address FROM processes JOIN listening_ports ON processes.pid = listening_ports.pid JOIN etc_protocols ON listening_ports.protocol = etc_protocols.number;",
			"interval": 300,
			"description": "List network processes + protocol"
		},
		"top_processes": {
			"query": "SELECT pid, name, uid, resident_size from processes order by resident_size desc limit 10;",
			"interval": 300,
			"description": "Select top 10 processes"
		},
		"most_active_processes": {
			"query": "select count(pid) as total, name, pid from processes group by name order by total desc limit 10;",
			"interval": 300,
			"description": "Return process count, name for the top 10 most active processes"
		},
		"select_non_persistent_processes": {
			"query": "SELECT name, path, pid FROM processes WHERE on_disk = 0;",
			"interval": 300,
			"description": ""
		},
		"osquery_info": {
			"query": "select * from osquery_info;",
			"interval": 300
		},
		"logged_in_users": {
			"query": "select * from logged_in_users;",
			"interval": 300,
			"description": ""
		},
		"startup_items": {
			"query": "select * from startup_items;",
			"interval": 300,
			"description": "# Select ALL startup items (to compare against a whitelist)"
		},
		"total_uptime_seconds": {
			"query": "select total_seconds from uptime;",
			"interval": 300
		},
		"interface_addresses": {
			"query": "select * from interface_addresses;",
			"interval": 3600
		}
	},
	"file_paths": {
		"homes": [
			"C:/Users/%/.ssh/%%"
		],
		"etc": [
			"C:/Windows/System32/drivers/etc/%%"
		],
		"destkop_links": [
			"C:/Users/%/Desktop/%.lnk"
		]
	},
	"exclude_paths": {},
	"decorators": {
		"load": [
			"SELECT uuid AS host_uuid FROM system_info;",
			"SELECT user AS username FROM logged_in_users ORDER BY time DESC LIMIT 1;",
			"SELECT ISO_8601 from time;",
			"SELECT address FROM interface_addresses WHERE mask NOT LIKE 'ffff%' AND address NOT LIKE '127.0.0.1' AND address NOT LIKE '169.254.%' AND address NOT LIKE '::%' AND address NOT LIKE '0.0.0.0' LIMIT 1;"
		]
	},
	"packs": {
		"osquery-monitoring": "C:/ProgramData/osquery/packs/osquery-monitoring.conf",
		"windows-attacks": "C:/ProgramData/osquery/packs/windows-attacks.conf"
	}
}
